# turbovanityurls/turbovanityurls.aur.install
post_upgrade() {
  /bin/systemctl restart turbovanityurls
}

pre_remove() {
  /bin/systemctl stop turbovanityurls
  /bin/systemctl disable turbovanityurls
}
