# notifiarr/notifiarr.aur.install
post_upgrade() {
  /bin/systemctl restart notifiarr
}

pre_remove() {
  /bin/systemctl stop notifiarr
  /bin/systemctl disable notifiarr
}
